# web3-auth-svc.service
[Unit]
Description=Web3 Auth Service (OPNet ML-DSA)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/web3-auth-svc-opnet
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
