# web3-auth-svc.service
[Unit]
Description=Web3 Auth Service (OPNet ML-DSA)
# Wait for libpam-web3-signing-host to resolve the host and generate a TLS
# cert with the matching SAN — without this, we may start with a cert
# whose SAN doesn't match the URL the user sees and they get a TLS warning.
After=network.target libpam-web3-signing-host.service
Wants=libpam-web3-signing-host.service

[Service]
Type=simple
ExecStart=/usr/bin/web3-auth-svc-opnet
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
